# ansible/roles/neo4j/files/neo4j.conf
description "neo4j upstart"

start on startup
stop on shutdown

exec sudo /opt/neo4j/neo4j-community-2.0.0/bin/neo4j start
